# config/systemd-rollback.service
[Unit]
Description=systemd-rollback


[Service]
Type=simple
ExecStart=/usr/local/bin/systemd-rollback
Restart=always

# Security reason not yet RND

SystemCallArchitectures=native
LockPersonality=yes
NoNewPrivileges=yes
PrivateTmp=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=CAP_SYS_PTRACE 
PrivateDevices=yes
PrivateUsers=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallFilter=~@privileged @resources
SystemCallFilter=@system-service


[Install]
WantedBy=multi-user.target